# rootfs-overlay/etc/systemd/system/seatd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/seatd -g video
